# ansible/templates/erp-backend.service.j2
[Unit]
Description=ERP System Spring Boot backend
After=network.target

[Service]
User={{ app_user }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/java -jar {{ app_dir }}/{{ jar_name }}
SuccessExitStatus=143
Restart=always
RestartSec=5

Environment=SPRING_PROFILES_ACTIVE=prod
Environment=SERVER_PORT={{ backend_port }}
Environment=DB_HOST={{ db_host }}
Environment=DB_PORT=3306
Environment=DB_NAME={{ db_name }}
Environment=DB_USER={{ db_user }}
Environment=DB_PASSWORD={{ db_password }}

[Install]
WantedBy=multi-user.target
